# scripts/sats/optuna-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Optuna Dashboard — launchd agent (persistent across reboots)

  Install (one-time):
    1. Edit REPO_PATH and VENV_PYTHON below to match your absolute paths.
    2. Copy to LaunchAgents:
         cp "/Volumes/Satechi Hub/warbird-pro/scripts/sats/optuna-dashboard.plist" \
            ~/Library/LaunchAgents/com.warbird.optuna-dashboard.plist
    3. Load the agent:
         launchctl load ~/Library/LaunchAgents/com.warbird.optuna-dashboard.plist
    4. Verify:
         launchctl list | grep warbird
         curl -s http://localhost:8080/ | head -1

  Manage:
    launchctl unload ~/Library/LaunchAgents/com.warbird.optuna-dashboard.plist  # stop
    launchctl load   ~/Library/LaunchAgents/com.warbird.optuna-dashboard.plist  # start

  Dashboard opens at: http://localhost:8080
  VS Code extension:
    - legacy SATS: data/sats_ps_optuna/study.db
    - other indicators: data/optuna/<indicator_key>/study.db
  For non-SATS studies, duplicate this plist and swap the sqlite path.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.warbird.optuna-dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Volumes/Satechi Hub/warbird-pro/.venv/bin/optuna-dashboard</string>
    <string>sqlite:////Volumes/Satechi Hub/warbird-pro/data/sats_ps_optuna/study.db</string>
    <string>--port</string>
    <string>8080</string>
    <string>--host</string>
    <string>127.0.0.1</string>
  </array>

  <!-- Restart if it crashes -->
  <key>KeepAlive</key>
  <true/>

  <!-- Start at login -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Working directory -->
  <key>WorkingDirectory</key>
  <string>/Volumes/Satechi Hub/warbird-pro</string>

  <!-- Log files -->
  <key>StandardOutPath</key>
  <string>/tmp/optuna-dashboard.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/optuna-dashboard.err</string>

  <!-- Wait for external drive to mount before starting -->
  <key>WatchPaths</key>
  <array>
    <string>/Volumes/Satechi Hub/warbird-pro/data/sats_ps_optuna</string>
  </array>
</dict>
</plist>
